# install/beatha.service
[Unit]
Description=Project Beatha - FPV Drone Recovery Tool
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/antoine/betaflightdebugger
ExecStart=/usr/bin/python3 src/backend/server.py
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/antoine/betaflightdebugger
Environment=LOG_LEVEL=INFO
StandardOutput=append:/var/log/beatha.log
StandardError=append:/var/log/beatha.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
